# scripts/systemd/cast-backup.service
[Unit]
Description=CAST backup (WAL-safe sqlite snapshot + runtime state)
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/opt/cast/app/scripts/backup.sh
